# install/100-usb.rules
#SHUTDOWN ON USB REMOVAL
#ANY USB:
#ACTION=="remove", SUBSYSTEM=="usb",RUN+="/sbin/shutdown -h now"
#
#USB STORAGE DEVICES:
KERNEL=="sd[a-z]*", ACTION=="remove", RUN+="/sbin/mpshutdown.sh"
